# curator4-7124533253.conf
[program:curator4-7124533253]
command=python3 /home/curator4/curator4-7124533253/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-7124533253/test.err.log
stdout_logfile=/home/curator4/curator4-7124533253/test.out.log

